# recipes-connectivity/alljoyn/alljoyn/alljoyn.service
[Unit]
Description=AllJoyn Daemon

[Service]
ExecStart=/usr/alljoyn/cpp/bin/alljoyn-daemon --internal

[Install]
WantedBy=multi-user.target
