# recipes-core/udev/udev-extraconf/bt-device.rules
# Set Bluetooth power up
ACTION=="add", SUBSYSTEM=="bluetooth", RUN+="/usr/bin/hciconfig %k up"
